# net-mail/davmail/files/initd
#!/sbin/runscript

CONFIG="/etc/${SVCNAME}.properties"
PIDFILE="/var/run/${SVCNAME}.pid"
SUFFIX="${SVCNAME#*.}"

if [[ "${SVCNAME}" == "${SUFFIX}" ]]; then
	SUFFIX=""
else
	SUFFIX=" (${SUFFIX})"
fi

depend() {
	need net
}

start() {
	ebegin "Starting DavMail${SUFFIX} gateway"

	if ! $(grep -Fl "davmail.server=true" "${CONFIG}" &> /dev/null); then
		eend 1 "Ensure that ${CONFIG} has the line davmail.server=true"
		return 1
	fi

	start-stop-daemon --start --user davmail --chdir /var/log/davmail \
		--make-pidfile --pidfile "${PIDFILE}" --background \
		--exec /usr/bin/davmail -- "${CONFIG}"

	eend $?
}

stop() {
	ebegin "Stopping DavMail${SUFFIX} gateway"
	start-stop-daemon --stop --user davmail --pidfile "${PIDFILE}"
	eend $?
}
