# www-servers/jboss/files/4.0.5/init.d/jboss-4
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/www-servers/jboss/files/3.2.5/init.d/jboss,v 1.2 2004/10/08 13:39:20 axxo Exp $

depend() {
    use net
}

start() {
	ebegin "Starting JBoss"
	if [ -z "$JBOSS_HOME" ]; then
		eerror "The JBOSS_HOME environment variable is not set."
		return 1
	fi
	if [ ! -d ${SERVER_HOME}/${JBOSS_CONF}/tmp ]; then
		mkdir ${SERVER_HOME}/${JBOSS_CONF}/tmp
		chown jboss:jboss ${SERVER_HOME}/${JBOSS_CONF}/tmp
	fi
	cd ${SERVER_HOME}/${JBOSS_CONF}/tmp
	if [ -z "${SUBIT}" ]; then
		eval ${CMD_START} >${JBOSS_CONSOLE} 2>&1 &
	else
		${SUBIT} "${CMD_START} >${JBOSS_CONSOLE} 2>&1 &" 
	fi
	eend $?
}

stop() {
	ebegin "Stopping JBoss"
	if [ -z "$SUBIT" ]; then
		${CMD_STOP}
	else
		${SUBIT} "$CMD_STOP"
	fi 
	eend $?
}
