# www-servers/jboss/files/4.0.5/init.d/jboss-4
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header$

depend() {
	use dns logger net
}

init_vars() {
	JBOSS_INT_OPTS="-Djboss.server.name=${JBOSS_SERVER_NAME}
		-Djboss.home.url=${JBOSS_HOME_URL}
		-Djboss.lib.url=${JBOSS_LIB_URL}
		-Djava.endorsed.dirs=\"${JBOSS_ENDORSED_DIRS}\"
		-Djboss.server.base.url=${JBOSS_SERVER_BASE_URL}
		-Djboss.server.home.url=${JBOSS_SERVER_HOME_URL}
		-Djboss.server.config.url=${JBOSS_SERVER_CONFIG_URL}
		-Djboss.server.root.deployment.filename=${JBOSS_SERVER_ROOT_DEPLOYMENT_FILENAME}
		-Djboss.server.data.dir=${JBOSS_SERVER_DATA_DIR}
		-Djboss.server.lib.url=${JBOSS_SERVER_LIB_URL}
		-Djboss.server.log.dir=${JBOSS_SERVER_LOG_DIR}
		-Djboss.server.temp.dir=${JBOSS_SERVER_TEMP_DIR}
		-Djboss.bind.address=${JBOSS_BIND_ADDRESS}
		-Djboss.server.exitonshutdown=${JBOSS_SERVER_EXIT_ON_SHUTDOWN}
		-Djboss.server.blockingshutdown=${JBOSS_SERVER_BLOCKING_SHUTDOWN}
		-Djboss.platform.mbeanserver=${JBOSS_PLATFORM_MBEAN_SERVER}
		-Djboss.server.requirejbossurlstreamhandlerfactory=${JBOSS_SERVER_REQUIRE_JBOSS_URL_STREAM_HANDLER_FACTORY}"
	if [ -n "${JBOSS_PATCH_URL}" ]; then
		JBOSS_INT_OPTS="${JBOSS_INT_OPTS} -Djboss.patch.url=${JBOSS_PATCH_URL}"
	fi
}

start()	{
	ebegin "Starting JBoss"

	MAX_FD_LIMIT=`ulimit -H -n`
	if [ $? -eq 0 ]; then
		if [ "${JBOSS_FD_LIMIT}" = "maximum" ]; then
			# use the system max
			JBOSS_FD_LIMIT="${MAX_FD_LIMIT}"
		fi
		ulimit -n ${JBOSS_FD_LIMIT}
		if [ $? -ne 0 ]; then
			ewarn "Could not set maximum file descriptor limit: ${JBOSS_FD_LIMIT}"
		fi
	else
		ewarn "Could not query system maximum file descriptor limit: ${MAX_FD_LIMIT}"
	fi

	init_vars
	start-stop-daemon --start --quiet --background --chuid ${JBOSS_USER}:${JBOSS_GROUP} \
		--make-pidfile --pidfile /var/run/jboss-4/${JBOSS_SERVER_NAME}/jboss-4.pid \
		--exec ${JAVA_HOME}/bin/java -- ${JAVA_OPTS} ${JBOSS_INT_OPTS} \
		-classpath ${JBOSS_STARTUP_CLASSPATH} org.jboss.Main
	eend $?
}

stop()	{
	ebegin "Stopping JBoss"
	start-stop-daemon --stop --quiet --make-pidfile --pidfile /var/run/${JBOSS_SERVER_NAME}/jboss-4.pid \
		--exec ${JAVA_HOME}/bin/java -- ${JAVA_OPTS} ${JBOSS_INT_OPTS} \
		-classpath ${JBOSS_SHUTDOWN_CLASSPATH} org.jboss.Shutdown ${JBOSS_SHUTDOWN_ARGS}
	eend $?
}
